# vw-backup.timer
[Unit]
Description=Vaultwarden Backup – täglich 03:00 Uhr

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
